# pyma.service
[Unit]
Description=PYMA
After=network.target

[Service]
ExecStart=/usr/bin/pyma --syslog

[Install]
WantedBy=multi-user.target
