# pyma.service
[Unit]
Description=PYMMA
After=network.target

[Service]
ExecStart=/usr/bin/pymma --syslog

[Install]
WantedBy=multi-user.target
